# rpi-camera.service
[Unit]
Description=RPi Camera MJPEG Server
After=network.target

[Service]
Type=simple
Restart=always
StartLimitInterval=0
StartLimitBurst=1440
RestartSec=60
ExecStart=/usr/local/bin/rpi-camera start
WorkingDirectory=/home/pi
User=pi
Group=pi
AmbientCapabilities=CAP_NET_BIND_SERVICE
# /run/rpi-camera holds the watchdog's restart marker. Owned by User=/Group=,
# preserved across `systemctl restart` (so the two-stage watchdog can tell a
# recovered restart from a still-wedged one) but wiped on reboot.
RuntimeDirectory=rpi-camera
RuntimeDirectoryPreserve=restart

[Install]
WantedBy=multi-user.target